# services/urbion-ai-utility-labelling.service
[Unit]
Description=Urbion AI Utility Labelling Service
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/projects/urbtec_utility_labelling
Environment="PATH=/home/ubuntu/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/ubuntu/.local/bin/uv run streamlit run __🏠_Home.py \
    --server.address 0.0.0.0 \
    --server.port 8501 \
    --server.fileWatcherType none \
    --server.runOnSave false \
    --browser.gatherUsageStats false
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
